# owt-ctrl/owt-net/deploy/deb/postinst
#!/bin/sh
set -e

ensure_oauth2_proxy_user() {
    if ! getent group oauth2-proxy >/dev/null 2>&1; then
        groupadd --system oauth2-proxy >/dev/null 2>&1 || addgroup --system oauth2-proxy >/dev/null 2>&1 || true
    fi
    if ! id -u oauth2-proxy >/dev/null 2>&1; then
        useradd --system --no-create-home --shell /usr/sbin/nologin --gid oauth2-proxy oauth2-proxy >/dev/null 2>&1 \
            || adduser --system --ingroup oauth2-proxy --no-create-home --disabled-login oauth2-proxy >/dev/null 2>&1 \
            || true
    fi
}

reset_owt_net_db() {
    db_path="/var/lib/owt-net/owt-net.db"
    db_dir="$(dirname "$db_path")"
    mkdir -p "$db_dir" >/dev/null 2>&1 || true
    rm -f "$db_path" >/dev/null 2>&1 || true
}

if [ "$1" = "configure" ]; then
    ensure_oauth2_proxy_user
    reset_owt_net_db

    if [ -f /etc/owt-net/oauth2-proxy.cfg ]; then
        chown root:oauth2-proxy /etc/owt-net/oauth2-proxy.cfg >/dev/null 2>&1 || true
        chmod 640 /etc/owt-net/oauth2-proxy.cfg >/dev/null 2>&1 || true
    fi
    if [ -f /etc/owt-net/google-allowed-emails.txt ]; then
        chown root:oauth2-proxy /etc/owt-net/google-allowed-emails.txt >/dev/null 2>&1 || true
        chmod 640 /etc/owt-net/google-allowed-emails.txt >/dev/null 2>&1 || true
    fi

    if [ -d /etc/nginx/sites-enabled ] && [ -f /etc/nginx/sites-available/owt-net.conf ] \
       && [ ! -e /etc/nginx/sites-enabled/owt-net.conf ]; then
        ln -s /etc/nginx/sites-available/owt-net.conf /etc/nginx/sites-enabled/owt-net.conf >/dev/null 2>&1 || true
    fi

    systemctl daemon-reload >/dev/null 2>&1 || true

    if systemctl is-enabled oauth2-proxy >/dev/null 2>&1; then
        systemctl restart oauth2-proxy >/dev/null 2>&1 || true
    else
        systemctl enable oauth2-proxy >/dev/null 2>&1 || true
        systemctl start oauth2-proxy >/dev/null 2>&1 || true
    fi

    if systemctl is-enabled owt-net >/dev/null 2>&1; then
        systemctl restart owt-net >/dev/null 2>&1 || true
    else
        systemctl enable owt-net >/dev/null 2>&1 || true
        systemctl start owt-net >/dev/null 2>&1 || true
    fi

    if command -v nginx >/dev/null 2>&1; then
        nginx -t >/dev/null 2>&1 && systemctl reload nginx >/dev/null 2>&1 || true
    fi
fi

exit 0
