# owt-ctrl/owt-net/deploy/deb/postinst
#!/bin/sh
set -e

reset_owt_net_db() {
    db_path="/var/lib/owt-net/owt-net.db"
    db_dir="$(dirname "$db_path")"
    mkdir -p "$db_dir" >/dev/null 2>&1 || true
    rm -f "$db_path" >/dev/null 2>&1 || true
}

if [ "$1" = "configure" ]; then
    reset_owt_net_db

    if [ -d /etc/nginx/sites-enabled ] && [ -f /etc/nginx/sites-available/owt-net.conf ] \
       && [ ! -e /etc/nginx/sites-enabled/owt-net.conf ]; then
        ln -s /etc/nginx/sites-available/owt-net.conf /etc/nginx/sites-enabled/owt-net.conf >/dev/null 2>&1 || true
    fi

    systemctl daemon-reload >/dev/null 2>&1 || true

    if systemctl is-enabled owt-net >/dev/null 2>&1; then
        systemctl restart owt-net >/dev/null 2>&1 || true
    else
        systemctl enable owt-net >/dev/null 2>&1 || true
        systemctl start owt-net >/dev/null 2>&1 || true
    fi

    if command -v nginx >/dev/null 2>&1; then
        nginx -t >/dev/null 2>&1 && systemctl reload nginx >/dev/null 2>&1 || true
    fi
fi

exit 0
